# supervisor.unison.conf
[program:unison]
command = unison -socket 5000
user = www-data
directory = %(ENV_UNISON_DIR)s
environment=HOME="/home/www-data",USER="www-data"
redirect_stderr = true
autorestart=true